# deploy/standalone/s6/s6-rc.d/mobazha-node/run
#!/bin/sh
exec 2>&1

ARGS="start -d /data"

if [ -n "$SAAS_PEER_ID" ]; then
  ARGS="$ARGS --httpproxytrustedpeer=$SAAS_PEER_ID"
fi

exec /usr/local/bin/mobazha $ARGS
