# deploy/standalone/s6/s6-rc.d/mobazha-node/run
#!/command/execlineb -P
/usr/local/bin/mobazha start -d /data
